# systemd/user/checkupdates.service
[Unit]
Description=pacman update checker for i3

[Service]
type=simple
ExecStart=/home/pit/.i3/checkupdates.sh
