# systemd/user/slideshow.service
[Unit]
Description=desktop wallpaper slideshow for i3

[Service]
type=simple
ExecStart=/home/pit/.i3/slideshow.sh
